# deploy/mmflux-whatsapp.service
[Unit]
Description=MMFlux WhatsApp Service
After=network.target mmflux.service

[Service]
Type=simple
User=mmflux
WorkingDirectory=/opt/mmflux/whatsapp-service
Environment=WA_PORT=3001
EnvironmentFile=/opt/mmflux/.env
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
